# mwms/consul/files/supervisor-client.conf
[program:consul-client]
command=/usr/local/sbin/consul agent -data-dir {{ consul_data_dir }} -config-dir /etc/consul
